# deploy/birdnet-collage.service
[Unit]
Description=birdnet-collage
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/birdnet-collage
ExecStart=/usr/bin/docker compose up -d --wait
ExecStop=/usr/bin/docker compose down
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target